# 20181006.04.bridge/config_files/files740nv5/etc/crontabs/root
30 4 * * *  [ $( cat /proc/uptime | cut -d '.' -f 1 ) -gt 3600 ] && reboot
*/3 * * * * /root/manage_hosts.sh add router.lan
*/3 * * * * /root/mqtt_control.sh start
